# etc/systemd/system/docker.service.d/proxy.conf
[Service]
#Environment="HTTP_PROXY=<PLACE_HOLDER>"
Environment="HTTPS_PROXY=<PLACE_HOLDER>"
Environment="NO_PROXY=localhost,127.0.0.1,.example.com"
